# pacific-shift-mcp-proxy/rootfs/etc/services.d/mcp-proxy/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Starting Pacific Shift MCP Proxy..."

exec python3 /usr/bin/start.py
